# setecna/rootfs/etc/services.d/setecna/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare REG_SYSTEM_ID
declare REG_USER
declare REG_PASSWORD
declare MQTT_HOST
declare MQTT_USER
declare MQTT_PASSWORD
declare ADV_INT
declare READONLY

## Get the 'systemID' key from the user config options.
REG_SYSTEM_ID="$(bashio::config 'systemID')"
REG_USER="$(bashio::config 'username')"
REG_PASSWORD="$(bashio::config 'password')"
MQTT_HOST="$(bashio::services mqtt 'host')"
MQTT_USER="$(bashio::services mqtt 'username')"
MQTT_PASSWORD="$(bashio::services mqtt 'password')"
ADV_INT="$(bashio::config 'adv_int')"
READONLY="$(bashio::config 'readonly')"

export REG_SYSTEM_ID
export REG_USER
export REG_PASSWORD
export MQTT_HOST
export MQTT_USER
export MQTT_PASSWORD
export ADV_INT
export READONLY

## Print the systemID the user supplied, defaults to ""
bashio::log.info "${systemID:=""}"

## Run your program
exec /usr/bin/app
